# output/jessie_rootfs/etc/systemd/system/cpu-corekeeper.service
[Unit]
Description=CPU corekeeper

[Service]
ExecStart=/usr/local/sbin/OrangePi_corekeeper.sh

[Install]
WantedBy=multi-user.target
